# configs/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/tmp/supervisord.
logfile_maxbytes=50MB
logfile_backups=10
loglevel = INFO
nodaemon=false
minfds=1024
minprocs=200
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/dev/shm/supervisor.sock
chmod = 0700
username = dummy
password = dummy

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock
username = dummy
password = dummy

[program:cron]
command=cron -f -L 15
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB

[include]
files=/etc/supervisor/conf.d/*.conf